# config/etc/services.d/fakes3/run
#!/usr/bin/with-contenv sh
fakes3 -r /fakes3/data -p 5670 --license LICENSE_KEY
